# thermitic/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:openssh]
command=/usr/sbin/sshd -D

[program:nginx]
command=/usr/sbin/nginx
