# thermitic/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:openssh]
command=/usr/sbin/sshd -D

[program:php5-fpm]
command=/usr/sbin/php5-fpm

[program:nginx]
command=/usr/sbin/nginx
